# server/irobot-server.service
[Unit]
Requires=network.target
After=network.target

[Service]
ExecStart=/usr/local/bin/irobot-server
Restart=always
